# email/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:queue-worker]
command=php artisan queue:work --tries=3 --timeout=120
redirect_stderr=true
stdout_logfile=/dev/stdout
